# contrib/package/node-id/files/etc/uci-defaults/node-id-defaults
#!/bin/sh
uci get system.@system[0].nodeid || uci set system.@system[0].nodeid=`ifconfig|md5sum |cut -f1 -d " "`
uci commit system
