# themes/piraten/root/etc/uci-defaults/luci-theme-piraten
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Piraten=/luci-static/piraten
        commit luci
EOF
	
