# contrib/package/piratenheartbeat/files/etc/uci-defaults/serverset
#!/bin/sh
uci get system.@pirate[0] || uci add system pirate
uci set system.@pirate[0].heartbeatserver=http://heartbeat.dd19.de/nodes/status/update
