# contrib/package/piratendefaults/files/etc/uci-defaults/piratendefaults
#!/bin/sh
uci batch <<-EOF
	set freifunk.community.homepage=http://berlin.piratenpartei.de
	commit freifunk
	set lucid.main.threadlimit=2
	commit lucid
EOF
	
